# tg-bot.service
[Unit]
Description=Telegram Group Bot
After=network.target

[Service]
Type=simple
User=your_username
WorkingDirectory=/path/to/tg-group-bot
Environment="PATH=/path/to/tg-group-bot/venv/bin"
ExecStart=/path/to/tg-group-bot/venv/bin/python3 /path/to/tg-group-bot/src/main.py
Restart=always
RestartSec=10
StandardOutput=append:/path/to/tg-group-bot/data/logs/bot.log
StandardError=append:/path/to/tg-group-bot/data/logs/bot_error.log

[Install]
WantedBy=multi-user.target